# config/systemd/certbot.service.d/onfailure.conf
# Source: config/systemd/certbot.service.d/onfailure.conf
# Drop-in for the stock certbot.service (Ubuntu 24.04 apt-managed).
# Adds OnFailure= wiring to surface renewal failures to journald.
# DO NOT add ExecStart or other overrides — this drop-in extends, not replaces.
# Installed to: /etc/systemd/system/certbot.service.d/onfailure.conf
# Applied by: scripts/bootstrap-edge.sh (operator-run)
[Unit]
OnFailure=certbot-renew-failure.service
